# elasticsearch.service
[Unit]
Description=elasticsearch
After=network.target
User=elasticsearch
Group=elasticsearch
RootDirectory=/var/lib/elasticsearch

[Service]
Type=simple
Restart=always
ExecStart=/usr/share/elasticsearch/bin/elasticsearch -c /var/lib/elasticsearch/elasticsearch.ini

[Install]
WantedBy=multi-user.target
